# storage/files/supervisor/worker.conf
[program:worker]
process_name=%(process_num)02d
command=php /var/www/html/artisan queue:listen --tries=1 --timeout=0
autostart=true
autorestart=false
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log